# claude_widget.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Claude Usage Widget — LaunchAgent (Autostart beim Login)
  ─────────────────────────────────────────────────────────
  WICHTIG: Passe den Pfad zu claude_widget.py unten an!

  Aktivieren:
    cp claude_widget.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/claude_widget.plist

  Deaktivieren:
    launchctl unload ~/Library/LaunchAgents/claude_widget.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chilli-mind.claude-widget</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/miniconda3/bin/python3</string>
        <string>/Users/kai/Desktop/Vibe-Code-Projekts/Claude-Usage-Widget/claude_widget.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- WICHTIG: Verbindet den Prozess mit dem macOS Window Server -->
    <!-- Ohne diesen Key sieht die Menüleisten-App keinen Display -->
    <key>SessionCreate</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kai/.cache/claude_widget.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kai/.cache/claude_widget.log</string>
</dict>
</plist>
